# platform/macos/com.nso.gc-controller.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  NSO GameCube Controller — macOS LaunchAgent

  Install to ~/Library/LaunchAgents/ to auto-start at login:
    cp com.nso.gc-controller.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.nso.gc-controller.plist

  Or enable from the app: Settings → Run at startup.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nso.gc-controller</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/NSO-GameCube-Controller-Pairing-App</string>
        <string>--minimized</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
