# platform/linux/99-gc-controller.rules
# Nintendo GameCube Controller Adapter (NSO)
# Allow access for HID and USB communication (users in "input" group)
SUBSYSTEM=="usb", ATTRS{idVendor}=="057e", ATTRS{idProduct}=="2073", GROUP="input", MODE="0660"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="057e", ATTRS{idProduct}=="2073", GROUP="input", MODE="0660"

# Auto-initialize the controller on plug-in so it emits standard HID reports.
# The gc-enabler binary sends the two USB bulk transfers and exits immediately.
# If gc-enabler is not installed, the rule is silently skipped.
#
# For headless emulation (BLE / non-SDL games), enable the user service:
#   systemctl --user enable gc-controller.service
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="057e", ATTRS{idProduct}=="2073", RUN+="/usr/local/bin/gc-enabler"

# Allow access to uinput for virtual gamepad creation (users in "input" group)
KERNEL=="uinput", GROUP="input", MODE="0660", OPTIONS+="static_node=uinput"
